# deploy/ngrok.service
[Unit]
Description=ngrok tunnel for hermes-mcp (alternative to cloudflared)
Documentation=https://ngrok.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Install as a systemd --user unit at ~/.config/systemd/user/ngrok.service.
# Run `ngrok config add-authtoken <token>` once before enabling this unit.
ExecStart=/usr/local/bin/ngrok http 8765 --domain=your-name.ngrok-free.app --log=stdout
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
